# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f
autostart=true
autorestart=true
stderr_logfile=/var/log/cron_err.log
stdout_logfile=/var/log/cron.log

; Настройки пользователя
user=root
